# system/LaunchAgents/local.claude.config-dir.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.claude.config-dir</string>
    <key>Description</key>
    <string>Configure Claude environment variables</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>launchctl setenv CLAUDE_CONFIG_DIR "$XDG_CONFIG_HOME/claude"</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
